# scripts/agentbridge-daemon.service
[Unit]
Description=AgentBridge (supervised-daemon)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory=/home/{{USER}}/.agentbridge
ExecStart=/home/{{USER}}/.agentbridge/watchdog.sh --all --web --agent
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=AGENT_BRIDGE_HOME=/home/{{USER}}/.agentbridge
Environment=ABMIND_HOME=/home/{{USER}}/.abmind
Environment=SUPERVISION=systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentbridge

[Install]
WantedBy=multi-user.target
