# scripts/agentbridge-watchdog.service
[Unit]
Description=AgentBridge Watchdog
After=network.target

[Service]
Type=simple
ExecStart=%h/.agentbridge/watchdog.sh --all --web --agent
Restart=on-failure
RestartSec=5
Environment=PATH=/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
